# evaluation/rules/windows/registry_event/sysmon_susp_lsass_dll_load.yml
filter: 'winlog.event_data.TargetObject: ("*\CurrentControlSet\Services\NTDS\DirectoryServiceExtPt*"
  OR "*\CurrentControlSet\Services\NTDS\LsaDbExtPt*")'
pre_detector:
  case_condition: directly
  id: b3503044-60ce-4bf4-bbcb-e3db98788823
  mitre:
  - attack.execution
  - attack.persistence
  - attack.t1177
  - attack.t1547.008
  severity: high
  title: DLL Load via LSASS
sigma_fields: true
description: Detects a method to load DLL via LSASS process using an undocumented
  Registry key