# evaluation/rules/windows/registry_event/sysmon_asep_reg_keys_modification.yml
filter: 'winlog.event_data.TargetObject: ("*\software\Microsoft\Windows\CurrentVersion\Run*"
  OR "*\software\Microsoft\Windows\CurrentVersion\RunOnce*" OR "*\software\Microsoft\Windows\CurrentVersion\RunOnceEx*"
  OR "*\software\Microsoft\Windows\CurrentVersion\RunServices*" OR "*\software\Microsoft\Windows\CurrentVersion\RunServicesOnce*"
  OR "*\software\Microsoft\Windows NT\CurrentVersion\Winlogon\Userinit*" OR "*\software\Microsoft\Windows
  NT\CurrentVersion\Winlogon\Shell*" OR "*\software\Microsoft\Windows NT\CurrentVersion\Windows\AppInit_DLLs*"
  OR "*\software\Wow6432Node\Microsoft\Windows NT\CurrentVersion\Windows\AppInit_DLLs*"
  OR "*\software\Microsoft\Windows NT\CurrentVersion\Windows\Load*" OR "*\software\Wow6432Node\Microsoft\Windows
  NT\CurrentVersion\Windows\Load*" OR "*\software\Microsoft\Windows NT\CurrentVersion\Windows\Run*"
  OR "*\software\Wow6432Node\Microsoft\Windows NT\CurrentVersion\Windows\Run*" OR
  "*\software\Microsoft\Windows\CurrentVersion\Explorer\User Shell Folders*")'
pre_detector:
  case_condition: directly
  id: 17f878b8-9968-4578-b814-c4217fc5768c
  mitre:
  - attack.persistence
  - attack.t1060
  - attack.t1547.001
  severity: medium
  title: Autorun Keys Modification
sigma_fields: true
description: Detects modification of autostart extensibility point (ASEP) in registry