# evaluation/rules/windows/registry_event/sysmon_apt_leviathan.yml
filter: 'winlog.event_data.TargetObject: "HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\Run\ntkd"'
pre_detector:
  case_condition: directly
  id: 70d43542-cd2d-483c-8f30-f16b436fd7db
  mitre:
  - attack.persistence
  - attack.t1060
  - attack.t1547.001
  severity: critical
  title: Leviathan Registry Key Activity
sigma_fields: true
description: Detects registry key used by Leviathan APT in Malaysian focused campaign